# user_service/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true

[program:grpc]
command=python manage.py grpcrunserver --dev
directory=/app
autostart=true
autorestart=true
